# deploy/systemd/money-retrain.timer
[Unit]
Description=Run Money nightly retrain

[Timer]
OnCalendar=*-*-* 02:15:00
AccuracySec=1m
Persistent=true
Unit=money-retrain.service

[Install]
WantedBy=timers.target
