# deploy/systemd/money-news.timer
[Unit]
Description=Run Money news feature refresh every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
