# deploy/systemd/money-api.service
[Unit]
Description=Money paper-trading FastAPI service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Update to the deployed repository path.
WorkingDirectory=/opt/money
# Place environment variables here or export them before starting the service.
EnvironmentFile=/etc/money/money.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/test -x /opt/money/.venv/bin/python
ExecStart=/opt/money/.venv/bin/python -u scripts/run_paper_api.py --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=45
SyslogIdentifier=money-api

[Install]
WantedBy=multi-user.target
