# deploy/systemd/money-retrain.timer
[Unit]
Description=Run Money nightly retrain

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true

[Install]
WantedBy=timers.target
